# curator4-1060408794.conf
[program:curator4-1060408794]
command=python3 /home/curator4/curator4-1060408794/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1060408794/test.err.log
stdout_logfile=/home/curator4/curator4-1060408794/test.out.log